# medatom/rootfs/src/etc/NetworkManager/dispatcher.d/03-lldpd
#!/bin/sh -e

interface="$1"
action="$2"

case $action in
	up)
		if /bin/systemctl is-active lldpd.service; then
			/bin/systemctl restart lldpd.service
		else
			/bin/systemctl start lldpd.service
		fi
		;;

	down)
		/bin/systemctl stop lldpd.service
		;;

	*)
		;;
esac
