# medcom-wide/rootfs/src/etc/systemd/system/lldpd.service
[Unit]
Description=LLDP Daemon

[Service]
ExecStart=/usr/sbin/lldpd -d -x -M 1
Sockets=lldpd.socket
StandardError=syslog

[Install]
WantedBy=multi-user.target
