# medcom-wide/rootfs/src/etc/systemd/system/nodm.service
[Unit]
Description=nodm Display Manager

[Service]
TimeoutStopSec=4sec
ExecStart=/usr/sbin/nodm

[Install]
WantedBy=graphical.target
