# medcom-wide/rootfs/src/etc/systemd/system/update-mime-db.service
[unit]
Description=Update MIME Type DB

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/lib/systemd/update-mime-db

[Install]
WantedBy=graphical.target
