# medcom-wide/rootfs/src/etc/systemd/system/var-lib.mount
[Unit]
Description=Library Directory
Before=local-fs.target

[Mount]
What=tmpfs
Where=/var/lib
Type=tmpfs
Options=mode=755

[Install]
WantedBy=local-fs.target
