# tamonten/tec/rootfs/src/etc/systemd/system/openssh.service
[Unit]
Description=OpenSSH Server

[Service]
ExecStart=/usr/sbin/sshd -D

[Install]
WantedBy=multi-user.target
