# medatom/rootfs/src/etc/systemd/system/lldpd.socket
[Unit]
Description=LLDP Daemon Socket

[Socket]
ListenStream=/var/run/lldpd.socket
