# medatom/rootfs/src/etc/systemd/system/slattach@.service
[Unit]
Description=SLIP on port /dev/%i

[Service]
ExecStart=/sbin/slattach /dev/%I -p slip -s 115200 -F -L

[Install]
WantedBy=multi-user.target
