# plutux/rootfs/src/etc/systemd/system/nodm.service
[Unit]
Description=nodm Display Manager

[Service]
ExecStart=/usr/sbin/nodm
TimeoutStopSec=4sec

[Install]
WantedBy=graphical.target
