# medatom-ng/rootfs/src/etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:/sbin/getty -L 38400 tty1 vt100
::respawn:/sbin/getty -L 38400 tty2 vt100
::respawn:/sbin/getty -L 38400 tty3 vt100
::respawn:/sbin/getty -L 38400 tty4 vt100
::respawn:/sbin/getty -L 38400 tty5 vt100
::respawn:/sbin/getty -L 38400 tty6 vt100
::respawn:/sbin/getty -L 115200 ttyS0 vt102
