# medatom-ng/rootfs/src/etc/systemd/system/usb-power.service
[Unit]
Description=USB Power Enable

[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/usb-power

[Install]
WantedBy=multi-user.target
